# bin/mirror_sync/offline_mirror.timer
[Unit]
Description=Keep offline mirror updated.

[Timer]
OnBootSec=1min
OnUnitActiveSec=4h
Unit=offline_mirror.service

[Install]
WantedBy=multi-user.target
