# bin/mirror_sync/offline_mirror.service
[Unit]
Description=Sync offline mirror.
 
[Service]
Type=simple
ExecStart=/usr/bin/bash -c "/mnt/mirror/projects/linux_install/bin/mirror_sync/mirror_sync.sh | tee -a /mnt/mirror/log.txt"

[Install]
WantedBy=multi-user.target
WantedBy=network.target
